# docker/supervisor/supervisord-mariadb.conf
[program:mariadb]
command=/usr/bin/mysqld_safe
numprocs=1
autostart=true
autorestart=true
priority=10
redirect_stderr=true
stdout_logfile=/var/log/mariadb/mariadb.log
stopwaitsecs=100
stopsignal=INT