# docker/supervisor/supervisord-apache2.conf
[program:apache2]
command=/usr/sbin/apache2ctl -DFOREGROUND
numprocs=1
autostart=true
autorestart=true
priority=50